# roles/znc/files/rc.d/znc
#!/bin/ksh

daemon="/usr/local/bin/znc"
daemon_flags="-d /var/znc"
daemon_user="_znc"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
